# crond/crontab
*	*	*	*	*	~/.syscfg/cron/dwm-statusbar.sh
